# deploy/aws-ec2/smartcontract.service.template
[Unit]
Description=Tokenized smartcontract Service

[Service]
EnvironmentFile={{configFilePath}}
ExecStart={{execPath}}

[Install]
WantedBy=multi-user.target
